# systemd/power.service.tmpl
[Unit]
Description=Power button monitor

[Service]
Type=simple
ExecStart=$TOP/scripts/power.py
Restart=always

[Install]
WantedBy=multi-user.target
Alias=power.service
